# supervisor-status/iml-supervisor-status.socket
[Unit]
Description=IML Supervisor Status Socket

[Socket]
ListenStream=/var/run/iml-supervisor-status.sock

[Install]
WantedBy=sockets.target
