# srcmap-reverse/iml-srcmap-reverse.socket
[Unit]
Description=Source Map Reverse Socket

[Socket]
ListenStream=/var/run/iml-srcmap-reverse.sock

[Install]
WantedBy=sockets.target
